# whatsapp/hcnews-send.timer
[Unit]
Description=HCNews WhatsApp Daily Send Timer
Documentation=file:///home/hc/Documentos/HCnews/whatsapp/PLAN.md

[Timer]
# Primary attempt at 6:00 AM
OnCalendar=*-*-* 06:00:00

# Retry every hour from 7:00 to 23:00 if not yet sent
OnCalendar=*-*-* 07:00:00
OnCalendar=*-*-* 08:00:00
OnCalendar=*-*-* 09:00:00
OnCalendar=*-*-* 10:00:00
OnCalendar=*-*-* 11:00:00
OnCalendar=*-*-* 12:00:00
OnCalendar=*-*-* 13:00:00
OnCalendar=*-*-* 14:00:00
OnCalendar=*-*-* 15:00:00
OnCalendar=*-*-* 16:00:00
OnCalendar=*-*-* 17:00:00
OnCalendar=*-*-* 18:00:00
OnCalendar=*-*-* 19:00:00
OnCalendar=*-*-* 20:00:00
OnCalendar=*-*-* 21:00:00
OnCalendar=*-*-* 22:00:00
OnCalendar=*-*-* 23:00:00

# If system was off at scheduled time, run when it comes back
Persistent=true

# Add some randomization to avoid exact-hour spikes
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
